# com.claude.fetch-usage.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.claude.fetch-usage</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/Users/andrei/.local/bin</string>
        <key>HOME</key>
        <string>/Users/andrei</string>
        <key>TERM</key>
        <string>xterm-256color</string>
    </dict>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/andrei/.claude/fetch-usage.sh</string>
    </array>
    <key>StartInterval</key>
    <integer>60</integer>
    <key>StandardOutPath</key>
    <string>/tmp/claude-fetch-usage.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/claude-fetch-usage.log</string>
</dict>
</plist>
